# crontab
# Ensure continuous jobs are running
*/5 * * * * jlocal /data/project/cluebot3/cluebot3/bigbrother.sh cluebot3 /data/project/cluebot3/cluebot3/start.sh
